# deploy/lila-http.socket
[Unit]
Description=Lichess HTTP socket

[Socket]
ListenStream=9001

[Install]
WantedBy=sockets.target
